# all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── ZooKeeper (start first) ───────────────────────────────────────────────────
[program:zookeeper]
command=/opt/kafka/bin/zookeeper-server-start.sh /etc/kafka/zookeeper.properties
priority=10
autostart=true
autorestart=true
startsecs=5
startretries=3
environment=KAFKA_HEAP_OPTS="-Xms256m -Xmx256m"
stdout_logfile=/var/log/supervisor/zookeeper.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/zookeeper.log
stderr_logfile_maxbytes=20MB

; ── Kafka (waits for ZooKeeper via wrapper script) ────────────────────────────
[program:kafka]
command=/usr/local/bin/start-kafka.sh
priority=20
autostart=true
autorestart=true
startsecs=20
startretries=5
environment=KAFKA_HEAP_OPTS="-Xms512m -Xmx512m"
stdout_logfile=/var/log/supervisor/kafka.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/kafka.log
stderr_logfile_maxbytes=20MB

; ── InfluxDB 2.7 ─────────────────────────────────────────────────────────────
[program:influxdb]
command=/usr/local/bin/influxd
    --bolt-path=/var/lib/influxdb2/influxd.bolt
    --engine-path=/var/lib/influxdb2/engine
    --storage-wal-fsync-delay=0
priority=20
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/influxdb.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/influxdb.log
stderr_logfile_maxbytes=20MB

; ── InfluxDB init (one-shot, runs after influxdb is up) ──────────────────────
[program:influxdb-init]
command=/usr/local/bin/influx-init.sh
priority=25
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/var/log/supervisor/influxdb-init.log
stderr_logfile=/var/log/supervisor/influxdb-init.log

; ── Glances dashboard data generator ─────────────────────────────────────────
[program:data-generator]
command=python3 /opt/data-generator/app.py
priority=26
autostart=true
autorestart=true
startsecs=3
startretries=10
environment=INFLUXDB_URL="http://localhost:8086",INFLUXDB_TOKEN="marketdata-token",INFLUXDB_ORG="test-org",INFLUXDB_BUCKET="glances",BACKFILL_DAYS="7",STEP_SECONDS="60",RUN_FOREVER="true",HOSTS="edge-01,edge-02",CONTAINERS="nginx,redis,postgres,worker",INTERFACES="eth0,eth1",DISKS="sda,sdb"
stdout_logfile=/var/log/supervisor/data-generator.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/data-generator.log
stderr_logfile_maxbytes=20MB

; ── Elasticsearch Prometheus dashboard generator/exporter ───────────────────
[program:elasticsearch-exporter-generator]
command=python3 /opt/data-generator/elasticsearch_exporter_generator.py --serve-only
priority=28
autostart=true
autorestart=true
startsecs=3
startretries=10
environment=ES_EXPORTER_LISTEN_HOST="0.0.0.0",ES_EXPORTER_LISTEN_PORT="9314",ES_BACKFILL_DAYS="7",ES_STEP_SECONDS="300",ES_SYNTHETIC_CLUSTER="synthetic-es-main",ES_SYNTHETIC_JOB="synthetic-elasticsearch-exporter",ES_NODES="es-hot-1,es-hot-2",PROMETHEUS_DATA_DIR="/prometheus",PROMTOOL_PATH="/usr/local/bin/promtool"
stdout_logfile=/var/log/supervisor/elasticsearch-exporter-generator.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/elasticsearch-exporter-generator.log
stderr_logfile_maxbytes=20MB

; ── Kafka Prometheus dashboard generator/exporter ───────────────────────────
[program:kafka-exporter]
command=python3 /opt/data-generator/kafka_exporter_generator.py --serve-only
priority=29
autostart=true
autorestart=true
startsecs=3
startretries=10
environment=KAFKA_EXPORTER_LISTEN_HOST="0.0.0.0",KAFKA_EXPORTER_LISTEN_PORT="9308",KAFKA_BACKFILL_DAYS="7",KAFKA_STEP_SECONDS="60",KAFKA_SYNTHETIC_NAMESPACE="kafka",KAFKA_BROKERS="broker-1,broker-2,broker-3",KAFKA_TOPICS="kraken-market-data,orders,payments",KAFKA_CONSUMER_GROUPS="analytics,etl,alerts",KAFKA_PARTITIONS_PER_TOPIC="6",KAFKA_REPLICATION_FACTOR="3",PROMETHEUS_DATA_DIR="/prometheus",PROMTOOL_PATH="/usr/local/bin/promtool"
stdout_logfile=/var/log/supervisor/kafka-exporter.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/kafka-exporter.log
stderr_logfile_maxbytes=10MB

; ── Elasticsearch 8.x ────────────────────────────────────────────────────────
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch
user=elasticsearch
priority=20
autostart=true
autorestart=true
startsecs=30
startretries=3
environment=ES_JAVA_OPTS="-Xms512m -Xmx512m"
stdout_logfile=/var/log/supervisor/elasticsearch.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/elasticsearch.log
stderr_logfile_maxbytes=20MB

; ── Prometheus ────────────────────────────────────────────────────────────────
[program:prometheus]
command=/usr/local/bin/start-prometheus.sh
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/prometheus.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/prometheus.log
stderr_logfile_maxbytes=20MB

; ── Loki ─────────────────────────────────────────────────────────────────────
[program:loki]
command=/usr/local/bin/loki -config.file=/etc/loki/config.yml
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/loki.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/loki.log
stderr_logfile_maxbytes=20MB

; ── node-exporter ────────────────────────────────────────────────────────────
[program:node-exporter]
command=/usr/local/bin/node_exporter
priority=30
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/node-exporter.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/node-exporter.log
stderr_logfile_maxbytes=10MB

; ── cadvisor (needs --privileged + host mounts for full metrics) ──────────────
[program:cadvisor]
command=/usr/local/bin/cadvisor
priority=30
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/cadvisor.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/cadvisor.log
stderr_logfile_maxbytes=10MB

; ── Grafana OSS ───────────────────────────────────────────────────────────────
[program:grafana]
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana --config=/etc/grafana/grafana.ini
priority=40
autostart=true
autorestart=true
environment=GF_PATHS_PROVISIONING="/etc/grafana/provisioning",GF_SECURITY_ADMIN_USER="admin",GF_SECURITY_ADMIN_PASSWORD="admin12345",GF_USERS_ALLOW_SIGN_UP="false",GRAFANA_PROMETHEUS_URL="http://localhost:9090",GRAFANA_INFLUXDB_URL="http://localhost:8086",GRAFANA_LOKI_URL="http://localhost:3100",GRAFANA_ELASTICSEARCH_URL="http://localhost:9200"
stdout_logfile=/var/log/supervisor/grafana.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/grafana.log
stderr_logfile_maxbytes=20MB

; ── Promtail ─────────────────────────────────────────────────────────────────
[program:promtail]
command=/usr/local/bin/promtail -config.file=/etc/promtail/config.yml
priority=50
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/promtail.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/promtail.log
stderr_logfile_maxbytes=10MB

; ── Filebeat ─────────────────────────────────────────────────────────────────
[program:filebeat]
command=/opt/filebeat/filebeat -e -c /etc/filebeat/filebeat.yml --strict.perms=false
priority=50
autostart=true
autorestart=true
startretries=10
stdout_logfile=/var/log/supervisor/filebeat.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/filebeat.log
stderr_logfile_maxbytes=10MB

; ── kafka-ui ─────────────────────────────────────────────────────────────────
[program:kafka-ui]
command=java -Xms256m -Xmx512m -jar /opt/kafka-ui.jar
priority=50
autostart=true
autorestart=true
startsecs=30
startretries=5
environment=KAFKA_CLUSTERS_0_NAME="local",KAFKA_CLUSTERS_0_BOOTSTRAPSERVERS="localhost:9092",KAFKA_CLUSTERS_0_ZOOKEEPER="localhost:2181",SERVER_PORT="8088"
stdout_logfile=/var/log/supervisor/kafka-ui.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/kafka-ui.log
stderr_logfile_maxbytes=20MB

; ── market-bridge (Python kraken-kafka-bridge) ────────────────────────────────
[program:market-bridge]
command=python3 -m kraken_kafka_bridge.main
priority=60
autostart=true
autorestart=true
startretries=20
environment=PYTHONPATH="/opt/market-bridge/src",KRAKEN_WS_URL="wss://ws.kraken.com/v2",KAFKA_BOOTSTRAP_SERVERS="localhost:9092",KAFKA_TOPIC="kraken-market-data",KRAKEN_CHANNEL="ticker",KRAKEN_SYMBOLS="BTC/USD,ETH/USD,XRP/USD",KRAKEN_BOOK_DEPTH="10",LOG_LEVEL="DEBUG",LOG_JSON_PATH="/var/log/market-bridge/market-bridge.jsonl",INFLUXDB_ENABLED="true",INFLUXDB_URL="http://localhost:8086",INFLUXDB_TOKEN="marketdata-token",INFLUXDB_ORG="test-org",INFLUXDB_BUCKET="marketdata"
stdout_logfile=/var/log/supervisor/market-bridge.log
stdout_logfile_maxbytes=20MB
stderr_logfile=/var/log/supervisor/market-bridge.log
stderr_logfile_maxbytes=20MB
